# private_dot_config/systemd/system/restart-calibre.path
[Unit]
Description=monitor calibre db
[Path]
PathModified=/home/raghu/calibre/library/metadata.db

[Install]
WantedBy=multi-user.target
