# private_dot_config/systemd/user/sync-calibre.path
[Path]
PathModified=/data/common/ebooks/metadata.db

[Install]
WantedBy=multi-user.target
